# ticket_lock_deploy/ticket-lock.service
[Unit]
Description=Ticket Lock Server (Autotask)
After=network.target

[Service]
Type=simple
User=ticketlock
WorkingDirectory=/opt/ticket-lock
Environment=TICKET_LOCK_API_KEY=CAMBIAR_ESTA_CLAVE
ExecStart=/opt/ticket-lock/venv/bin/gunicorn -w 2 -b 127.0.0.1:5001 ticket_lock_server:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
